# launchd/com.playground.auto-news-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.playground.auto-news-bot</string>

    <key>ProgramArguments</key>
    <array>
      <string>/usr/bin/python3</string>
      <string>/Users/sergejsergeev/Documents/Playground/auto_news_bot/main.py</string>
      <string>--config</string>
      <string>/Users/sergejsergeev/Documents/Playground/auto_news_bot/config.json</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/sergejsergeev/Documents/Playground/auto_news_bot</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PYTHONUNBUFFERED</key>
      <string>1</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/sergejsergeev/Documents/Playground/auto_news_bot/data/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/sergejsergeev/Documents/Playground/auto_news_bot/data/stderr.log</string>
  </dict>
</plist>
